# ops/roles/tonight/files/tonight_oauth2.service
[Unit]
Description=Job that runs the oauth2 proxy for tonight

[Service]
WorkingDirectory=/var/www/tonight/bitly_oauth2_proxy
ExecStart=/var/www/tonight/bitly_oauth2_proxy/oauth2_proxy --config=oauth2_proxy.cfg

[Install]
WantedBy=multi-user.target
